# usm_peleng.cron
* * * * *    root    docker-compose -f /path/to/userside/docker-compose.yml run --rm usm_peleng perl usm_peleng.pl &> /dev/null
